# images/dib/elements/hotstack-sonic-vs/extra-data.d/container-files/rebootbackend.conf
# Dummy rebootbackend program to satisfy start.sh
# This service doesn't exist in SONiC-VS but start.sh tries to start it
[program:rebootbackend]
command=/bin/sleep infinity
autostart=false
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog
