# devsetup/hotstack-os/systemd/hotstack-os-placement.service
# Copyright Red Hat, Inc.
# All Rights Reserved.
#
# Licensed under the Apache License, Version 2.0 (the "License"); you may
# not use this file except in compliance with the License. You may obtain
# a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS, WITHOUT
# WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied. See the
# License for the specific language governing permissions and limitations
# under the License.

[Unit]
Description=HotStack-OS Placement Resource Tracking Service
Documentation=https://github.com/openstack-k8s-operators/hotstack
After=hotstack-os-mariadb.service hotstack-os-keystone.service hotstack-os-dnsmasq.service hotstack-os-haproxy.service
Requires=hotstack-os-mariadb.service hotstack-os-keystone.service hotstack-os-dnsmasq.service hotstack-os-haproxy.service
PartOf=hotstack-os.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30
SuccessExitStatus=143

ExecStartPre=-/usr/bin/podman stop -t 2 hotstack-os-placement
ExecStartPre=-/usr/bin/podman rm -f hotstack-os-placement
ExecStart=/usr/bin/podman run \
    --pull=never \
    --name hotstack-os-placement \
    --hostname placement \
    --network hotstack-os \
    --ip __PLACEMENT_IP__ \
    --env PLACEMENT_DB_PASSWORD=__DB_PASSWORD__ \
    --env SERVICE_PASSWORD=__SERVICE_PASSWORD__ \
    --env KEYSTONE_ADMIN_PASSWORD=__KEYSTONE_ADMIN_PASSWORD__ \
    --env REGION_NAME=__REGION_NAME__ \
    --volume __HOTSTACK_DATA_DIR__/runtime/config/placement/placement.conf:/etc/placement/placement.conf:ro,z \
    localhost/hotstack-os-placement:latest

# Health check: CONTAINER INTERVAL TIMEOUT RETRIES START_PERIOD COMMAND
ExecStartPost=-/usr/local/bin/hotstack-healthcheck.sh hotstack-os-placement 10 5 5 5 curl -f http://localhost:8778/
ExecStop=/usr/bin/podman stop -t 30 hotstack-os-placement
ExecStopPost=-/usr/bin/podman rm -f hotstack-os-placement

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=hotstack-os.target
